# stonehouse/install/stonehouse.conf
# Start stonehouse on Vagrant mounted

start on vagrant-mounted

/vagrant/install/stonehouse.sh