# application/install/stonehouse.conf
# Start stonehouse on Vagrant mounted

start on vagrant-mounted

/vagrant/stonehouse-bundle/stonehouse.sh